# user/yubikey-touch-detector.service.d/override.conf
[Unit]
After=gpg-agent.service gpg-agent.socket gpg-agent-ssh.socket
PartOf=org.gnome.Shell.target
