# system/x.target
[Unit]
Description=VT-less X display servers
Before=display-manager.service

[Install]
RequiredBy=display-manager.service
