# user/spotifyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/spotifyd --no-daemon --cache-path %C/spotifyd --device-name '%l'
