# user/syncthing.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/syncthing serve --no-browser --no-restart --logflags=0 $SYNCTHING_ARGS
Environment=SYNCTHING_ARGS=""
EnvironmentFile=-%E/syncthing/syncthing.env
