# user/cpu.slice
[Unit]
Description=CPU-intensive tasks

[Slice]
CPUAccounting=true
MemoryAccounting=true
