# user/exit-session@.target
[Unit]
Description=Exit the graphical session %I
DefaultDependencies=no
