# user/fetch-bcachefs.timer
[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
