# user/graphic@.target
.include need-session@.target

[Unit]
Description=Graphical session %I
Requires=overall.target
AllowIsolate=yes
