# user/jupyter.service
[Service]
ExecStart=/usr/bin/jupyter notebook --no-browser --log-level=INFO
Restart=on-failure

[Install]
WantedBy=default.target
