# user/kdeinit-module@.service.in
.include need-session@.service

[Unit]
Description=KDE module @MODULE@ for session %I
Before=kdeinit@%i.target
Wants=kdeinit@%i.target
After=dbus.service kdeinit@%i.service

[Service]
Type=oneshot
ExecStart=/usr/bin/dbus-send --session --type=method_call --dest=org.kde.kded /kded org.kde.kded.loadModule string:@MODULE@

[Install]
WantedBy=graphic@.target
