# user/modprobed-db.timer
[Unit]
Description=Periodically collect data about loaded modules

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
