# user/overall.target
[Unit]
Description=Persistent non-session-bound user services
Requires=basic.target
AllowIsolate=yes

[Install]
WantedBy=default.target
