# user/overall.target
[Unit]
Description=Non-graphical part of user session
Requires=basic.target
AllowIsolate=yes

[Install]
WantedBy=default.target
