# user/primary-xserver.target

[Unit]
Description=Primary X server set
