# user/shairport.service
[Unit]
Description=Shairport AirTunes receiver
Requires=sound.target
After=sound.target

[Service]
EnvironmentFile=%h/.config/shairport
ExecStart=/usr/bin/shairport $SHAIRPORT_ARGS
Restart=always

[Install]
WantedBy=basic.target
