# user/shairport.service
[Unit]
Description=AirPort Express emulator
Requires=pulseaudio.socket

[Service]
EnvironmentFile=%h/.config/shairport
ExecStart=/usr/bin/shairport-sync -o pulse $SHAIRPORT_ARGS
Restart=on-failure

[Install]
WantedBy=default.target
