# user/udiskie.service
[Service]
Type=simple
ExecStart=/usr/bin/udiskie --verbose --config %E/udiskie/config.yml
